# supervisord-owncloud.conf
[program:owncloud]
command=/root/scripts/start-owncloud.sh
numprocs=1
autostart=true
autorestart=false
startsecs=0
priority=100